# meshcore-bot.service
[Unit]
Description=MeshCore LLM Bot (Jeff)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/jeff
Environment="PATH=/home/YOUR_USERNAME/jeff/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/YOUR_USERNAME/jeff/.env
ExecStart=/home/YOUR_USERNAME/jeff/venv/bin/python3 /home/YOUR_USERNAME/jeff/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
